# .config/systemd/user/gitlab-projects.mount.template
[Unit]
Description=Github Project Folder Bind Mounts
Before=local-fs.target umount.target

[Mount]
What=@@BIND_MOUNT@@
Where=%h/gitlab
Type=fuse.bindfs

[Install]
WantedBy=default.target
